# .config/systemd/user/weather.service
[Unit]
Description=Run weather data fetch
ConditionNetworkOnline=true

[Service]
Type=oneshot
ExecStart=%h/.local/bin/weather.py getdata
